# Impotfiles/supervisor/conf.d/newyolo.ini
[program:newyolo]
command=gunicorn --bind 0.0.0.0:8080 newyolo.wsgi
directory=/root/newyolo/
autostart=true
autorestart=unexpected
user=root
stdout_logfile=/etc/supervisor/var/log/newyolo.log
stderr_logfile=/etc/supervisor/var/log/newyolo.log

